# deploy/github-to-signal.service
[Unit]
Description=GitHub webhook to Signal notifications
After=signal-cli-bot.service
Requires=signal-cli-bot.service

[Service]
Type=exec
ExecStart=/usr/local/bin/github-to-signal
WorkingDirectory=/home/tai/code/foss/github-to-signal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
