# deploy/signal-cli-bot.service
[Unit]
Description=signal-cli daemon for github-to-signal bot
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/signal-cli -a +18134455076 daemon --http 127.0.0.1:8081 --no-receive-stdout
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
